# Build/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid
childlogdir=/var/log/supervisord/
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error

[program:erupe]
stdout_logfile=/Erupe/logs/stdout.log
stdout_logfile_maxbytes=0
redirect_stdout=true
stderr_logfile=/Erupe/logs/stderr.log
stderr_logfile_maxbytes=0
redirect_stderr=true

command=/Erupe/erupe-ce
autostart=true
autorestart=true
